# scripts/client/crontab
EDAKIA_PATH=/home/tux/project/dakia
* * * * * cd $EDAKIA_PATH; sh scripts/export.sh >/dev/null 2>&1
0 */2 * * *  find $EDAKIA_PATH/files/* \! -name "*.log" -type f -mmin +360 -delete >/dev/null 2>&1
