# scripts/client/crontab
* * * * * cd /home/tux/project/dakia; sh scripts/export.sh >/dev/null 2>&1
0 */2 * * *  find /home/tux/project/dakia/files/* \! -name "*.log" -type f -mmin +360 -delete >/dev/null 2>&1
